# supervisord.conf
[supervisord]
nodaemon=true

[program:cntlm]
command=/usr/sbin/cntlm -U cntlm -f

[program:redsocks]
command=/usr/sbin/redsocks
